# scripts/dropbox
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
    need net
}

start() {
    ebegin "Odpalanie Dropboxa"
    su -c /home/asgavar/.dropbox-dist/dropboxd asgavar
    eend $?
}
